# raspi3/openstack/etc/uci-defaults/xxx_config
uci set wireless.@wifi-device[0].disabled="0"
uci set wireless.@wifi-device[0].country="US"

uci set wireless.@wifi-iface[0].disabled="0"
uci set wireless.@wifi-iface[0].ssid="BH-Network"
uci set wireless.@wifi-iface[0].key="Sli07ih01!"
uci set wireless.@wifi-iface[0].mode="sta"
uci set wireless.@wifi-iface[0].encryption="psk2"
uci commit wireless

uci set network.lan.ipaddr="192.168.5.1"
uci commit network

uci set system.@system[0].hostname=beige
uci commit system

uci set memcached.@memcached.memory=256
uci commit memcached