# raspi3/files/etc/uci-defaults/xxx_config
uci set wireless.@wifi-device[0].disabled="0"
uci set wireless.@wifi-device[1].disabled="0"
uci set wireless.@wifi-iface[1].disabled="0"
uci set wireless.@wifi-iface[1].ssid="BH-Network"
uci set wireless.@wifi-iface[1].key="Sli07ih01!"
uci set wireless.@wifi-iface[1].mode="sta"
uci set wireless.@wifi-iface[1].encryption="psk2"
uci commit wireless

uci set system.@system[0].hostname=kore
uci commit system